# data/Dockerfiles/logger/supervisord.conf
[supervisord]
nodaemon=true

[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground  --no-caps
redirect_stderr=true
autostart=true
stdout_syslog=true

[program:postfix]
command=/opt/logger.sh
autorestart=true

[supervisorctl]
serverurl=unix:///var/tmp/supervisord.sock

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
